# rootfs/etc/s6-overlay/s6-rc.d/ttyd/run
#!/bin/bash
exec s6-setuidgid claude \
    /usr/local/bin/ttyd \
    --port 7681 \
    --credential "${TTYD_USERNAME}:${TTYD_PASSWORD}" \
    /bin/bash -l
