# rootfs/etc/s6-overlay/s6-rc.d/ttyd/run
#!/bin/bash
exec s6-setuidgid claude \
    s6-env HOME=/home/claude \
    /usr/local/bin/ttyd \
    -W \
    --port 7681 \
    --ping-interval 30 \
    --cwd /home/claude/workspace \
    --credential "${TTYD_USERNAME}:${TTYD_PASSWORD}" \
    -t rendererType=canvas \
    /bin/bash -l
