# rootfs/etc/s6-overlay/s6-rc.d/tailscaled/run
#!/bin/bash

# Tailscale is optional — if TS_AUTHKEY is not set, exit cleanly
# so s6 does not keep restarting the daemon unnecessarily.
if [ -z "$TS_AUTHKEY" ]; then
    exec s6-svc -Od .
fi

mkdir -p /home/hatchpod/.tailscale
chown hatchpod:hatchpod /home/hatchpod/.tailscale
mkdir -p /var/run/tailscale

# Auto-detect TUN device availability for transparent networking
# Check for the character device only — do NOT read from it (reads block/fail)
if [ -c /dev/net/tun ]; then
    # Kernel TUN mode — transparent routing, no proxy needed
    echo "tailscaled: using kernel TUN mode (transparent networking)" >&2
    exec tailscaled \
        --statedir=/home/hatchpod/.tailscale \
        --socket=/var/run/tailscale/tailscaled.sock
else
    # Userspace networking fallback — apps need SOCKS5 proxy for Tailscale peers only
    echo "tailscaled: /dev/net/tun not available, falling back to userspace networking" >&2

    # Write proxy env vars scoped to Tailscale traffic only.
    # We do NOT set ALL_PROXY / HTTP_PROXY / HTTPS_PROXY because that would
    # redirect ALL internet traffic through the Tailscale SOCKS5 proxy and
    # break general connectivity. Users who need proxy-based access to
    # Tailscale peers can source this file or set these vars manually.
    cat > /etc/profile.d/tailscale-proxy.sh <<'EOF'
# Tailscale userspace-networking SOCKS5 proxy (auto-generated)
# NOT exported by default — setting ALL_PROXY/HTTP_PROXY globally breaks
# internet connectivity for non-Tailscale destinations.
# To use: export ALL_PROXY=socks5h://localhost:1055
# Or use per-command: ALL_PROXY=socks5h://localhost:1055 curl http://tailscale-peer/...
TAILSCALE_PROXY=socks5h://localhost:1055
EOF

    exec tailscaled \
        --tun=userspace-networking \
        --statedir=/home/hatchpod/.tailscale \
        --socket=/var/run/tailscale/tailscaled.sock
fi
