# rootfs/etc/s6-overlay/s6-rc.d/api/run
#!/bin/bash
exec s6-setuidgid claude \
    s6-env HOME=/home/claude \
    s6-env HOST=0.0.0.0 \
    s6-env BROWSE_ROOT=/home/claude/workspace \
    s6-env DEFAULT_CWD=/home/claude/workspace \
    node /opt/api-server/dist/index.js
