# rootfs/etc/s6-overlay/s6-rc.d/tailscaled-up/up
/etc/s6-overlay/scripts/tailscaled-up.sh
